# demo/slurm/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/slurm/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:munged]
command=/usr/sbin/munged --foreground --key-file=/etc/munge/munge.key
user=munge
autostart=true
autorestart=true
stdout_logfile=/var/log/slurm/munge.log
stderr_logfile=/var/log/slurm/munge.log
priority=10

[program:slurmctld]
command=/usr/sbin/slurmctld -D
user=slurm
autostart=true
autorestart=true
stdout_logfile=/var/log/slurm/slurmctld.log
stderr_logfile=/var/log/slurm/slurmctld.log
startsecs=3
priority=20

[program:slurmd]
command=/usr/sbin/slurmd -D -N slurm-node
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/slurm/slurmd.log
stderr_logfile=/var/log/slurm/slurmd.log
startsecs=5
priority=30
